# scripts/fintrack-backup.timer
[Unit]
Description=Run FinTrack Pro database backup daily

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
